# config/etc/cron.d/adafruit_sync
0-59/5 * * * * pi /bin/sh /opt/energyconsole/bin/send_adafruit.sh
